# systemd/system/netns-helper@.target
[Unit]
Description=Network namespace %i target
Requires=netns-helper-postup@%i.service
AllowIsolate=no

[Install]
WantedBy=multi-user.target